# modules/nat_instance/snat/snat.service
[Unit]
Description = NAT Instance

[Service]
ExecStart = /opt/snat/snat.sh
Type = oneshot

[Install]
WantedBy = multi-user.target